# roles/000-base-011-haproxy_exporter/templates/haproxy_exporter.service.j2
[Unit]
Description=Prometheus Server
Documentation=https://prometheus.io/docs/introduction/overview/
After=network.target

[Service]
Restart=on-failure
ExecStart={{ haproxy_exporter_bin }}

[Install]
WantedBy=multi-user.target