# roles/000-base-012-consul_exporter/templates/consul_exporter.service.j2
[Unit]
Description=Prometheus Server
Documentation=https://prometheus.io/docs/introduction/overview/
After=network.target

[Service]
Restart=on-failure
ExecStart={{ consul_exporter_bin }}

[Install]
WantedBy=multi-user.target