# roles/000-base-010-node_exporter/templates/node_exporter.service.j2
[Unit]
Description=Prometheus Server
Documentation=https://prometheus.io/docs/introduction/overview/
After=network.target

[Service]
Restart=on-failure
ExecStart={{ node_exporter_bin }}

[Install]
WantedBy=multi-user.target